# companion_scripts/cloudflared-hls.service
[Unit]
Description=Cloudflared Quick Tunnel for HLS Camera Stream
Documentation=https://developers.cloudflare.com/cloudflare-one/networks/connectors/cloudflare-tunnel/
After=network-online.target
Wants=network-online.target

# Start before camera-stream so the tunnel is ready when the stream registers
Before=camera-stream.service

[Service]
Type=simple
User=julius
# Tunnel local HLS server (port 8080) to a random *.trycloudflare.com URL.
# --metrics 127.0.0.1:33843 exposes a fixed metrics port so camera_stream_service.py
# can auto-detect the public URL via http://127.0.0.1:33843/quicktunnel
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:8080 --metrics 127.0.0.1:33843
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared-hls

[Install]
WantedBy=multi-user.target
